# modules/power-management/systemd/lid-close-handler.service
[Unit]
Description=Lid Close Handler for Hyprland Dream
Documentation=https://github.com/dreamcoder08/hyprland-dream
After=graphical-session.target
Wants=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
User=%i
Environment=DISPLAY=:0
Environment=WAYLAND_DISPLAY=wayland-0
ExecStart=/usr/bin/bash %h/.config/hyprland-dream/modules/power-management/scripts/lid-close-handler.sh
Restart=on-failure
RestartSec=2
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=%h/.config %h/.local/share

# Power management
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=5

[Install]
WantedBy=graphical-session.target 